# dot_config/autostart/com.dropbox.Client.desktop
[Desktop Entry]
Type=Application
Name=Dropbox
Exec=flatpak run com.dropbox.Client
X-GNOME-Autostart-enabled=true
